# deployment/podman-systemd-service-units/rating-update.container
[Container]
AutoUpdate=registry
ContainerName=rating-update
Environment=TZ="America/Chicago"
Image=ghcr.io/rnevils/rating-update:dev-latest
Pod=rating-updates.pod